# extra/dist/tarantool.service.in
# It's not recommended to modify this file in-place, because it will be
# overwritten during package upgrades. If you want to customize there're
# number of ways:

# Recommended way:
# 1) Use "/etc/sysconfig/tarantool" or "/etc/default/tarantool" -
# They're supported by our start-stop utility - dist.lua

# Usual way for RPM-based distros
# 2) Create a file "/etc/systemd/system/tarantool.service",
# containing
#   .include /usr/lib/systemd/system/tarantool.service
#   # Here're your changes
#
# For example, if you want to change CONF_DIR create
# "/etc/systemd/system/tarantool.service" containing:
#   .include /usr/lib/systemd/system/tarantool.service
#   [Service]
#   Environment=CONF_DIR=/etc/tarantool/instances.other
# This will override the settings appearing below

[Unit]
Description=Tarantool instances
After=network.target
Documentation=man:tarantool(1)

[Service]
Type=forking
OOMScoreAdjust=-1000

ExecStart=@SYSV_INITD_PATH@ start
ExecStop=@SYSV_INITD_PATH@ stop

TimeoutSec=300
