# FreeBSD/databases/tarantool/files/tarantool.in
#!/bin/sh
#
# $FreeBSD: databases/tarantool/files/tarantool.in 308708 2012-12-11 20:00:36Z rm $

# PROVIDE: tarantool
# REQUIRE: LOGIN
# KEYWORD: shutdown
#
# tarantool_enable="YES"
# tarantool_config=""
#

. /etc/rc.subr

name="tarantool"
rcvar=tarantool_enable

load_rc_config "$name"

tarantool_enable=${tarantool_enable:-"NO"}
tarantool_config=${tarantool_config:-"%%ETCDIR%%/$name.cfg"}

command="%%PREFIX%%/bin/tarantool_box"
command_args="--daemonize --config ${tarantool_config}"
pidfile="/var/run/$name.pid"

run_rc_command "$1"
