# debian/tarantool-common.cron.hourly
#!/bin/sh

test -x /usr/bin/tarantool_box || exit 0
test -x /usr/sbin/tarantool_snapshot_rotate || exit 0

exec /usr/sbin/tarantool_snapshot_rotate
